# packaging/systemd/wattpost.service
[Unit]
Description=WattPost solar-monitor daemon
Documentation=https://github.com/ritualnorth/wattpost
# Bluetooth must be up before we try to discover the BT-2 dongle.
# network-online.target lets the cloud sync + ntfy / Discord / SMTP
# transports reach the outside world on first poll.
After=bluetooth.service network-online.target
Wants=bluetooth.service network-online.target
# Auto-rollback watchdog: fires when the StartLimit below is hit.
OnFailure=wattpost-rollback.service

[Service]
Type=simple
User=wattpost
Group=wattpost
# bluetooth group gives access to BlueZ DBus without root.
SupplementaryGroups=bluetooth

WorkingDirectory=/opt/wattpost
Environment=PYTHONUNBUFFERED=1
# Single source of truth for runtime paths so the install script and
# pi-gen stage only have to set them in one place.
Environment=WATTPOST_CONFIG=/etc/wattpost/config.yaml
Environment=WATTPOST_DB=/var/lib/wattpost/solar-monitor.db
Environment=WATTPOST_PORT=80

# Allow binding privileged ports without running as root. AmbientCapabilities
# grants the daemon CAP_NET_BIND_SERVICE so the default :80 dashboard works as
# User=wattpost. Docker / manual installs override WATTPOST_PORT to something
# unprivileged via /etc/wattpost/port.env (installed by install.sh).
#
# CapabilityBoundingSet is deliberately NOT pinned to CAP_NET_BIND_SERVICE:
# capping it there strips CAP_SETUID/SETGID from the daemon's children, so
# `sudo` (setuid root) can't escalate ("unable to change to root gid") and the
# firewall / SSH / update helpers all fail. Left at the default set so sudo
# works; the proper fix is a root helper service off sudo (tracked separately).
EnvironmentFile=-/etc/wattpost/port.env
AmbientCapabilities=CAP_NET_BIND_SERVICE

ExecStart=/opt/wattpost/venv/bin/solar-monitor --log-level INFO serve \
    --config ${WATTPOST_CONFIG} \
    --db ${WATTPOST_DB} \
    --interval 60 \
    --port ${WATTPOST_PORT} \
    --host 0.0.0.0

# /api/system/restart uses os.execv to re-exec the daemon. Restart=
# on-failure means systemd also brings us back if the process crashes
# for any reason that *doesn't* end in an execv.
Restart=on-failure
RestartSec=5
TimeoutStopSec=15

# Auto-rollback watchdog (#36 Slice 3). If the daemon fails to
# start 3 times within 60s, systemd gives up on the unit AND fires
# the OnFailure unit below — which rolls /opt/wattpost back to the
# previous slot and starts the daemon there. Pairs with the
# inactive-slot health probe in wattpost-update: the probe catches
# import errors / API breakage, this watchdog catches "boots fine
# in sandbox but crashes against real DB / real BLE / real config".
StartLimitIntervalSec=60
StartLimitBurst=3

# Hardening — minimum capabilities, read-only filesystem except for
# the state + log paths we explicitly need.
#
# The daemon must invoke `sudo` for the root helpers the sudoers entries
# pin: wattpost-update / wattpost-rollback (Update/Rollback buttons) and
# wattpost-netctl (the firewall + SSH toggles, cloud#15). sudo is setuid
# root, so it only works while NoNewPrivileges stays OFF.
#
# CAUTION (learned the hard way): NoNewPrivileges is forced ON *implicitly*
# by several hardening directives — ProtectKernelTunables, ProtectKernelModules,
# ProtectControlGroups, RestrictNamespaces, LockPersonality all imply it. We
# had those set, so despite never setting NoNewPrivileges= explicitly, sudo was
# silently blocked ("sudo: unable to change to root gid") — breaking the
# firewall/SSH toggles on every install. They're removed here so sudo works.
# The proper fix is to move the privileged ops off sudo into a root-run helper
# service, letting the daemon keep full hardening (tracked separately). Until
# then we keep only the directives that do NOT imply NoNewPrivileges:
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
# config.yaml gets rewritten in place by the setup wizard + alerts
# editor, so /etc/wattpost has to be writable too. The dnsmasq-shared.d
# entry lets the daemon drop/remove its captive-portal DNS-hijack file
# (install.sh group-owns the dir to wattpost, but ProtectSystem=strict
# still binds it read-only without this). Leading `-` keeps the unit
# starting on Docker installs where NetworkManager — and that dir —
# don't exist.
ReadWritePaths=/var/lib/wattpost /etc/wattpost -/etc/NetworkManager/dnsmasq-shared.d

# journald captures stdout/stderr; the in-memory ring buffer also keeps
# the last 500 lines for Settings → Diagnostics.
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
